# deploy/central/systemd/asterisk.service
[Unit]
Description=Asterisk PBX (DIDStorage)
After=network.target postgresql.service redis-server.service didapi.service
Wants=didapi.service

[Service]
Type=simple
User=asterisk
Group=asterisk
ExecStart=/usr/sbin/asterisk -f -U asterisk -G asterisk -vvvg -C /etc/asterisk/asterisk.conf
ExecReload=/usr/sbin/asterisk -rx 'core reload'
Restart=always
RestartSec=4
LimitCORE=infinity
LimitNOFILE=65536
LimitNPROC=infinity

# Ambient cap so asterisk (running as non-root) can bind :5060 / RTP ports.
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
