# deploy/central/systemd/sip-capture.service
[Unit]
Description=Persistent SIP + RTP packet capture for CDR sip-trace + call-quality analysis
After=network.target

[Service]
Type=simple
User=root
# Capture SIP signalling (5060/5061) AND Asterisk's default RTP media port
# range (10000-20000) so the call-quality analyser has RTP streams to work
# from. Rolling pcap files, daily rotation, keep 7 files (~7 days).
# %%Y / %%m / %%d become Y/m/d at runtime (systemd needs %% to escape).
#
# The RTP range matches Asterisk's out-of-the-box rtp.conf. If you widen it
# there, widen it here too, or the quality report will show one-way audio
# for the missing side.
ExecStartPre=/bin/mkdir -p /var/lib/didstorage/sip-traces
ExecStart=/usr/bin/tcpdump -i any -n -s 0 -W 7 -G 86400 \
  -w /var/lib/didstorage/sip-traces/sip-%%Y%%m%%d.pcap \
  -Z root \
  "udp and (port 5060 or port 5061 or portrange 10000-20000) and not ip6"
Restart=always
RestartSec=4

[Install]
WantedBy=multi-user.target
